# crontab
*/1 * * * * date >> ~/nicetomeetyou/nba.log ; python3 ~/nicetomeetyou/crawler.py >> ~/nicetomeetyou/nba.log
